# templates/etc/systemd/system/nftables.service.d/override.conf.j2
# {{ ansible_managed }}
[Service]
ExecStart=
ExecStart=/usr/sbin/nft -f /etc/nftables/nftables.conf
ExecReload=
ExecReload=/usr/sbin/nft -f /etc/nftables/nftables.conf
